# udev/99-usb-notify.rules
# USB block device connect/disconnect notifications
SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", ENV{ID_BUS}=="usb", ACTION=="add", RUN+="/usr/local/bin/usb-notify.sh"
SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", ENV{ID_BUS}=="usb", ACTION=="remove", RUN+="/usr/local/bin/usb-notify.sh"
